# deploy/systemd/mirror-manager.service
[Unit]
Description=apt-mirror manager (web dashboard for repos + disk usage)
Documentation=file:///opt/apt/mirror-manager/README
After=network.target

[Service]
Type=simple
# PRODUCTION: run as a DEDICATED NON-ROOT user — scripts/setup-mirror-manager.sh creates
# that user plus a narrow sudoers grant and a drop-in that overrides User= below. This base
# unit runs as root only as a fallback (it writes /etc/apt/mirror.list + /opt/apt/keys and
# starts apt-mirror.service); a bug in the HTTP handler while root is an RCE-to-root path,
# so prefer the non-root setup.
User=root
WorkingDirectory=/opt/apt/mirror-manager
ExecStart=/usr/bin/python3 /opt/apt/mirror-manager/mirror_manager.py
Restart=on-failure
RestartSec=3
# Bind to localhost only; expose via SSH tunnel or the nginx auth proxy.
Environment=MM_LISTEN_HOST=127.0.0.1
Environment=MM_LISTEN_PORT=8080
# Set a shared secret if you proxy it; clients pass the X-MM-Token header.
# Environment=MM_TOKEN=change-me
# Environment=MM_BUDGET_BYTES=1700000000000

# Sandboxing — confines a bug in the HTTP handler; writes limited to /etc/apt + /opt/apt.
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/etc/apt /opt/apt
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
RestrictNamespaces=yes
RestrictSUIDSGID=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
LockPersonality=yes
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
# Aggressive extras — enable once validated in your environment:
# CapabilityBoundingSet=
# MemoryDenyWriteExecute=yes

[Install]
WantedBy=multi-user.target
