# deploy/systemd/mirror-manager-ldap-auth.service
[Unit]
Description=LDAPS auth backend for nginx auth_request (apt-manager dashboard)
Documentation=file:///opt/apt/mirror-manager/README
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Runs as the dashboard's apt-manager user and CO-LOCATED with the dashboard, so it can
# read the UI-managed auth files (/opt/apt/manager/{ldap.json,htpasswd}) and verify local
# users with `htpasswd -vb`. It still only talks to LDAPS + answers localhost.
# (Separate-host proxy without the dashboard? Use DynamicUser=yes and configure LDAP purely
# via the Environment= lines below — but then the Access tab can't manage it.)
User=apt-manager
Group=apt-manager
WorkingDirectory=/opt/apt/mirror-manager
ExecStart=/usr/bin/python3 /opt/apt/mirror-manager/ldap_auth.py
Restart=on-failure
RestartSec=3

Environment=LA_LISTEN_HOST=127.0.0.1
Environment=LA_LISTEN_PORT=8889
# UI-managed auth files (written by the dashboard; take precedence over the env below).
Environment=LDAP_CONF=/opt/apt/manager/ldap.json
Environment=HTPASSWD=/opt/apt/manager/htpasswd
# ---- fallback directory config (used only when ldap.json is absent) ----
Environment=LDAP_URI=ldaps://dc.example.com:636
Environment=LDAP_BIND_MODE=direct
Environment=LDAP_USER_DN_TEMPLATE={user}@example.com
# Optional: restrict to a group, verify the LDAPS cert against a CA, etc.
# Environment=LDAP_REQUIRED_GROUP=CN=apt-admins,OU=Groups,DC=example,DC=com
# Environment=LDAP_CA=/etc/ssl/certs/example-root-ca.pem
# Environment=LDAP_TLS_REQCERT=demand
# ---- search-mode alternative (service account) ----
# Environment=LDAP_BIND_MODE=search
# Environment=LDAP_BIND_DN=CN=svc-apt,OU=Service,DC=example,DC=com
# Environment=LDAP_BIND_PW=change-me
# Environment=LDAP_BASE_DN=DC=example,DC=com
# Environment=LDAP_USER_FILTER=(sAMAccountName={user})

NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
